# scripts/play.service
[Unit]
Description=PLAY - Sistema de Reprodução de Vídeos
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/PLAY
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

